# docker/cron/crontab
* * * * * cd /srv/basis_member && php ctos app:hello-cron >> /var/log/cron/basis_member.log 2>&1