# docker/cron/crontab
* * * * * cd /srv/data_exchange && php ctos app:hello-cron >> /var/log/data_exchange/data_exchange.log 2>&1
0 0 * * 0 /usr/sbin/logrotate /etc/logrotate.d/custom_logs